# launchd/com.apolitical.assistant.backfill.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.apolitical.assistant.backfill</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/npm</string>
        <string>run</string>
        <string>backfill</string>
        <string>--</string>
        <string>--from=yesterday</string>
    </array>

    <key>WorkingDirectory</key>
    <string>APOLITICAL_ASSISTANT_PATH</string>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>0</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <key>StandardOutPath</key>
    <string>APOLITICAL_ASSISTANT_PATH/logs/launchd/backfill.log</string>

    <key>StandardErrorPath</key>
    <string>APOLITICAL_ASSISTANT_PATH/logs/launchd/backfill.error.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
